# System/local.d/backlight.start
#!/bin/sh

echo 1333 > /sys/class/backlight/intel_backlight/brightness
